# flaskapp.service
[Unit]
Description=MARBEFES BBT Database Flask Application
Documentation=https://marbefes.eu
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=exec
User=razinka
Group=razinka
WorkingDirectory=/home/razinka/OneDrive/HORIZON_EUROPE/MARBEFES/EMODNET_PyDeck
Environment=PATH=/home/razinka/OneDrive/HORIZON_EUROPE/MARBEFES/EMODNET_PyDeck/venv/bin
Environment=VIRTUAL_ENV=/home/razinka/OneDrive/HORIZON_EUROPE/MARBEFES/EMODNET_PyDeck/venv
Environment=FLASK_APP=app.py
Environment=FLASK_ENV=production
Environment=FLASK_DEBUG=0
Environment=HOST=0.0.0.0
Environment=PORT=5000
Environment=WMS_BASE_URL=https://ows.emodnet-seabedhabitats.eu/geoserver/emodnet_view/wms
Environment=HELCOM_WMS_BASE_URL=https://maps.helcom.fi/arcgis/services/MADS/Pressures/MapServer/WMSServer
Environment=WMS_TIMEOUT=30
Environment=VECTOR_DATA_DIR=data/vector
Environment=LOG_LEVEL=INFO
Environment=LOG_FILE=logs/app.log
Environment=SECRET_KEY=your-production-secret-key-change-this
ExecStartPre=/home/razinka/OneDrive/HORIZON_EUROPE/MARBEFES/EMODNET_PyDeck/venv/bin/python -c "import flask, requests; print('Dependencies OK')"
ExecStart=/home/razinka/OneDrive/HORIZON_EUROPE/MARBEFES/EMODNET_PyDeck/venv/bin/python /home/razinka/OneDrive/HORIZON_EUROPE/MARBEFES/EMODNET_PyDeck/run_flask.py
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s TERM $MAINPID
KillMode=mixed
TimeoutStartSec=60
TimeoutStopSec=30
RestartSec=10
Restart=on-failure
RestartPreventExitStatus=1

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/razinka/OneDrive/HORIZON_EUROPE/MARBEFES/EMODNET_PyDeck/logs
ReadWritePaths=/home/razinka/OneDrive/HORIZON_EUROPE/MARBEFES/EMODNET_PyDeck/data
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=marbefes-flask

[Install]
WantedBy=multi-user.target