# marbefes-bbt.service
[Unit]
Description=MARBEFES BBT Database - Marine Biodiversity and Ecosystem Functioning Database
After=network.target

[Service]
Type=notify
User=razinka
Group=razinka
WorkingDirectory=/home/razinka/OneDrive/HORIZON_EUROPE/MARBEFES/EMODNET_PyDeck
Environment="PATH=/home/razinka/anaconda3/bin:/usr/local/bin:/usr/bin"
EnvironmentFile=/home/razinka/OneDrive/HORIZON_EUROPE/MARBEFES/EMODNET_PyDeck/.env

# Gunicorn configuration
ExecStart=/home/razinka/anaconda3/bin/gunicorn -c gunicorn_config.py app:app
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s TERM $MAINPID

# Process management
Restart=on-failure
RestartSec=5s
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Security settings
PrivateTmp=true
NoNewPrivileges=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=marbefes-bbt

[Install]
WantedBy=multi-user.target
