# marbefes-bbt-database.service
[Unit]
Description=MARBEFES BBT Database - Marine Biodiversity and Ecosystem Functioning Database
Documentation=https://marbefes.eu
After=network.target
Requires=network.target

[Service]
Type=notify
User=razinka
Group=razinka
WorkingDirectory=/home/razinka/OneDrive/HORIZON_EUROPE/MARBEFES/EMODNET_PyDeck
Environment=PATH=/home/razinka/OneDrive/HORIZON_EUROPE/MARBEFES/EMODNET_PyDeck/venv/bin
EnvironmentFile=/home/razinka/OneDrive/HORIZON_EUROPE/MARBEFES/EMODNET_PyDeck/.env.production

# Gunicorn command
ExecStart=/home/razinka/OneDrive/HORIZON_EUROPE/MARBEFES/EMODNET_PyDeck/venv/bin/gunicorn \
    --config /home/razinka/OneDrive/HORIZON_EUROPE/MARBEFES/EMODNET_PyDeck/gunicorn.conf.py \
    --env FLASK_ENV=production \
    app:app

ExecReload=/bin/kill -s HUP $MAINPID

# Restart configuration
Restart=always
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=3

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/home/razinka/OneDrive/HORIZON_EUROPE/MARBEFES/EMODNET_PyDeck/logs
ReadWritePaths=/home/razinka/OneDrive/HORIZON_EUROPE/MARBEFES/EMODNET_PyDeck/data

# Resource limits
LimitNOFILE=65535
TimeoutStartSec=30
TimeoutStopSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=marbefes-bbt-database

[Install]
WantedBy=multi-user.target